# output-root/etc/hotplug.d/iface/01-wportal
#when lan ip changed, need to restart wportal, rmmod here .

[ -f /tmp/boot_done ] || exit 0
[ "$INTERFACE" == "lan" ] && wportalctrl -c && echo "stop" > /tmp/wportal/status