# output-root/etc/hotplug.d/iface/26-client_mgmt
[ "$INTERFACE" == "lan" -a "$ACTION" == "ifup" ] && {
    /etc/init.d/client_mgmt restart
}
