# output-root/etc/hotplug.d/iface/98-cloud-brd
[ -f /tmp/boot_done ] || exit 0
[ "$ACTION" == "ifup" ] && cloud_reconnect
