# output-root/etc/hotplug.d/button/40-wireless
state="off"
changed=0

obtain_state() {
	local disabled disabled_all band
	local dev="$1"
	config_get disabled $1 disabled
	config_get disabled_all $1 disabled_all
	config_get band $1 band

	if [ "$ACTION" = "released" -a "$BUTTON" = "wps" -a $SEEN -ge 2 ]; then
		if [ "$disabled" = "off" ]; then
			state="on"
			ledcli WIFI2G_OFF
			ledcli WIFI5G_OFF
			ledcli WPS_OFF
		else
			if [ "$band" = "2g" -a "$disabled_all" = "off" ]; then
				ledcli WIFI2G_ON
			fi
			if [ "$band" = "5g" -a "$disabled_all" = "off" ]; then
				ledcli WIFI5G_ON
			fi
		fi
	fi
}

do_toggle() {
    local dev="$1"

    if [ "$ACTION" = "released" -a "$BUTTON" = "wps" -a $SEEN -ge 2 ]; then
        uci set wireless.$dev.disabled=$state
        wifi_commit
	changed=1
    fi
}


toggle_wireless(){
	config_load wireless
	config_foreach obtain_state wifi-device
	config_foreach do_toggle wifi-device
	local switch="on"
	[ $state = "on" ] && switch="off"
	[ $changed -eq 1 ] && logwrite $$ wireless WIRELESS_BUTTON_SET $switch
}	

is_factory_mode(){
	config_load factory
	config_get enable  factorymode enable
	echo "factory_mode:" $enable >/dev/console
	[ $enable = "yes" ] && echo "true"  || echo "false"
}

dut_is_cal_when_boot(){
	[ -f /tmp/dut_is_not_cal ] && echo "false"  || echo "true"
}

check(){
    if [ $(is_factory_mode) = "true" -a $(dut_is_cal_when_boot) = "false" ]; then
        echo "wifi button will lose function now!" >/dev/console
        touch /tmp/button_wifi_check
    elif [ $SEEN -ge 2 ]; then
        toggle_wireless
        
        config_load sysmode
        config_load onemesh 
        config_get sysmode sysmode mode "router"
        config_get onemesh_enable onemesh enable "on"	
        
        local support_onemesh=$(uci get profile.@onemesh[0].onemesh_support -c "/etc/profile.d" -q)
        if [ "$support_onemesh" == "yes" -a "$sysmode" = "router" -a "$onemesh_enable" = "on" ]; then
            lua -e 'require("luci.controller.admin.onemesh").sync_wifi_all()'
        fi

        wifi
    fi
}

[ "$ACTION" = "released" -a "$BUTTON" = "wps" ] && check

 

