# output-root/etc/init.d/br-filter
#!/bin/sh /etc/rc.common
# Copyright (C) 2008-2010 OpenWrt.org


start() {
	echo ath01 ath11 > /proc/bridge_filter/device_list
}
