# output-root/etc/init.d/sys-led-end
#!/bin/sh /etc/rc.common

START=99
SERVICE_DAEMONIZE=1
cycle()
{
	sleep 1
	ledcli STATUS_ON
}

start()
{
	cycle & 
}
